# init/snapraid-btrfs-runner.timer
[Unit]
Description=Run snapraid-btrfs-runner

[Timer]
OnCalendar=daily
RandomizedDelaySec=5h
Persistent=true

[Install]
WantedBy=timers.target
